# docker/encountergen/conf/supervisord/conf.d/ionic.conf
[program:ionicserver]
command = ionic serve --no-open
directory = /encountergenux