# docker/monstergen/conf/supervisord/conf.d/server.conf
[program:encountergen]
command = gunicorn main:app -b localhost:5000 --reload
directory = /encountergen